# my-stream.service
[Unit]
Description=my-stream file watcher
After=network.target

[Service]
Type=simple
User=deploy                          # cambia con il tuo utente sulla VPS
WorkingDirectory=/var/www/my-stream  # cambia con il percorso del tuo progetto
ExecStart=/usr/bin/node watch.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
